# debian/prerm
#!/bin/bash
# prerm script for indi-mqtt

systemctl stop indi-mqtt.service
systemctl disable indi-mqtt.service
